# crontab
* * * * * curl http://localhost/backup.php >> /var/www/html/backups/backup.log
#